# platform/broadcom/sonic-platform-modules-cel/debian/platform-modules-blackstone.postinst
depmod -a
systemctl enable platform-modules-blackstone.service
systemctl start platform-modules-blackstone.service

/usr/local/bin/platform_api_mgnt.sh install