# platform/innovium/sonic-platform-modules-cel/debian/platform-modules-mrvl-cpo1.postinst
depmod -a
systemctl enable platform-modules-mrvl-cpo1.service
systemctl start platform-modules-mrvl-cpo1.service
